# roles/janitor/templates/janitor.cron.j2
#!/bin/sh

export DOCKER_API_VERSION='1.18'
/opt/janitor-{{ janitor_version}}/janitor
